# ServoBlaster/udev_scripts/20-servoblaster.rules
SUBSYSTEM=="module", DEVPATH=="/module/servoblaster", RUN+="/lib/udev/servoblaster"
